# deploy/appliance/examples/orchestrator.service.d/50-helm-boundary.conf
# Boundary Enforcement Profile appliance-vector-01 (tier enforce).
# Compiled by `helm-ai-kernel boundary profile compile`. Do not edit: posture drift fails closed.
# Sealed topology: orchestrator.service may reach only the HELM gateway.
[Service]
IPAddressDeny=any
IPAddressAllow=127.0.0.1
