# reference_packs/boundary-profile-v1/artifacts/systemd/helm-gateway.service.d/50-helm-boundary.conf
# Boundary Enforcement Profile appliance-vector-01 (tier enforce).
# Compiled by `helm-ai-kernel boundary profile compile`. Do not edit: posture drift fails closed.
[Service]
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CPUQuota=50%
MemoryMax=536870912
TasksMax=128
DevicePolicy=closed
DeviceAllow=/dev/null rw
